# ansible/roles/vmhost/templates/docker-template.timer.j2
[Unit]
Description={{ item.cmd }}

[Timer]
OnCalendar={{ item.schedule }}

[Install]
WantedBy=timers.target
