# ansible/roles/vmhost/templates/docker-template.service.j2
[Unit]
Description={{ item.cmd }}

[Service]
ExecStart={{ config_path }}{{ run_cli_tool }} {{ item.cmd }}
